# dist/librenms-scheduler.service
[Unit]
Description=LibreNMS Scheduler

[Service]
Type=oneshot
StandardError=journal
StandardOutput=journal
WorkingDirectory=/opt/librenms/
ExecStart=php artisan schedule:run --no-ansi --no-interaction
Environment="PATH=/opt/librenms/venv/bin:/usr/local/bin:/usr/bin:/bin"
User=librenms
Group=librenms
